# scripts/vps/systemd/indemnite-indexnow.timer
[Unit]
Description=Daily IndexNow ping — fires ~30min after the daily rebuild
Documentation=https://www.indexnow.org/

[Timer]
# Daily rebuild fires at 16:00 UTC (±30 min jitter). We fire at 16:30 UTC
# with our own jitter so the rebuild has time to push + Cloudflare to
# deploy. The script itself adds another 180s sleep, so by the time we
# actually fetch the sitemap, we're ~30-35 minutes after rebuild push.
# Slot distinct from calc (09:30 UTC) and concursoja (11:00 UTC) per
# footprint-hygiene staggering.
OnCalendar=*-*-* 16:30:00 UTC
RandomizedDelaySec=30min
Persistent=true
Unit=indemnite-indexnow.service

[Install]
WantedBy=timers.target
