# scripts/vps/systemd/indemnite-indexnow.service
[Unit]
Description=Calcul Indemnité IndexNow ping — submit today's newly-surfaced URLs to Bing/Yandex/DDG
Wants=network-online.target
After=network-online.target indemnite-rebuild.service

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/indemnite
ExecStart=/home/ubuntu/indemnite/scripts/vps/indexnow-ping.sh
StandardOutput=journal
StandardError=journal
Restart=no
# Script waits 180s for Cloudflare rebuild + has retries; cap total at 10min
TimeoutStartSec=600
# Non-fatal: IndexNow being down should not flag a "service failed"
SuccessExitStatus=0 1
