# scripts/vps/systemd/indemnite-scrape@.service
[Unit]
Description=Calcul Indemnité scraper: %i
Wants=network-online.target
After=network-online.target postgresql.service
# On a failed run, arm a transient retry timer (+6h, +30h). A single transient
# blip would otherwise go stale until the next weekly run and trip the 8d health
# alert (cf calculify incident 2026-06-15).
OnFailure=indemnite-scrape-retry@%i.timer

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/indemnite/pipeline
EnvironmentFile=/home/ubuntu/indemnite/.env.local
ExecStart=/home/ubuntu/indemnite/pipeline/.venv/bin/python -m pipeline.cli scrape %i
StandardOutput=journal
StandardError=journal
# Don't restart automatically; the timer schedules retries.
Restart=no
# Kill after 5 minutes — well above the seconds we actually need.
TimeoutStartSec=300
