# scripts/vps/systemd/renov-health.service
[Unit]
Description=MaRénovAide scrape pipeline health check
After=network-online.target postgresql.service

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/renov/pipeline
EnvironmentFile=/home/ubuntu/renov/.env.local
ExecStart=/home/ubuntu/renov/pipeline/.venv/bin/python -m pipeline.cli health
StandardOutput=journal
StandardError=journal
# Non-zero exit means issues; that's expected behavior, not service failure.
SuccessExitStatus=0 1
