# scripts/vps/systemd/renov-rebuild.service
[Unit]
Description=MaRénovAide daily rebuild — empty commit to unlock scheduled drip-publishes
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/renov
EnvironmentFile=/home/ubuntu/renov/.env.local
ExecStart=/home/ubuntu/renov/scripts/vps/daily-rebuild.sh
StandardOutput=journal
StandardError=journal
Restart=no
TimeoutStartSec=180
